# petkit-ble-bridge/rootfs/etc/services.d/petkit-ble-bridge/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare address
declare mqtt
declare mqtt_broker
declare mqtt_port
declare mqtt_user
declare mqtt_password
declare logging_level

## Get the keys from the user config options.
address=$(bashio::config 'address')
mqtt_broker=$(bashio::config 'mqtt_broker')
mqtt_port=$(bashio::config 'mqtt_port')
mqtt_user=$(bashio::config 'mqtt_user')
mqtt_password=$(bashio::config 'mqtt_password')

## Run your program
exec /etc/petkit-ble-bridge/run \
  --address "${address}" \
  --mqtt \
  --mqtt_broker "${mqtt_broker}" \
  --mqtt_port ${mqtt_port} \
  --mqtt_user "${mqtt_user}" \
  --mqtt_password "${mqtt_password}" \
  --logging_level "INFO"